# supervisord/supervisord.conf
[unix_http_server]
file=/tmp/supervisor.sock   ; path to your socket file
;chmod=0777
;chown=webscan:webscan

[inet_http_server]
port=1222
username=david
password=123456

[supervisord]
logfile=/var/log/supervisord/supervisord.log ; supervisord log file
logfile_maxbytes=50MB       ; maximum size of logfile before rotation
logfile_backups=10          ; number of backed up logfiles
loglevel=info               ; info, debug, warn, trace
pidfile=/var/run/supervisord.pid ; pidfile location
nodaemon=false              ; run supervisord as a daemon
minfds=1024                 ; number of startup file descriptors
minprocs=200                ; number of process descriptors
user=root                  ; default user
childlogdir=/var/log/supervisord/            ; where child log files will live

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock ; use unix:// schem for a unix sockets.

;=============================================celery flower==================================================
[program:flower]
command=/usr/local/bin/flower --adress=0.0.0.0 --port=5555 --broker=amqp://rose:rosecelery@localhost:5672/celery
autostart=false
autorestart=true
user=rose
directory=/home/rose/MStore/scripts
stdout_logfile_maxbytes = 50MB
stdoiut_logfile_backups = 20
stdout_logfile = /var/log/flower.log

;=============================================celery main worker==================================================
[program:main_worker]
; Set full path to celery program if using virtualenv
command=celery worker -A tasks -l info -Q main -n main

environment=PYTHONPATH=/home/rose/MStore

directory=/home/rose/MStore/scripts
user=rose
numprocs=1
stdout_logfile=/var/log/celery/worker_main.log
stderr_logfile=/var/log/celery/worker_main.log
autostart=false
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 60

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=999

;=============================================celery image worker==================================================
[program:image_worker]
; Set full path to celery program if using virtualenv
command=celery worker -A tasks -l info -Q download -n image-%(program_name)s-%(process_num)d

process_name = %(program_name)s-%(process_num)d

environment=PYTHONPATH=/home/rose/MStore

directory=/home/rose/MStore/scripts
user=rose
numprocs=2
stdout_logfile=/var/log/celery/worker_download.log
stderr_logfile=/var/log/celery/worker_download.log
autostart=false
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 600

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=998

;=============================================celery crawl worker==================================================
[program:crawl_worker]
; Set full path to celery program if using virtualenv
command=celery worker -A tasks -l info -Q default -n crawl-%(program_name)s-%(process_num)d

environment=PYTHONPATH=/home/rose/MStore

directory=/home/rose/MStore/scripts
user=rose
numprocs=1
stdout_logfile=/var/log/celery/worker_download.log
stderr_logfile=/var/log/celery/worker_download.log
autostart=false
autorestart=true
startsecs=10

; Need to wait for currently executing tasks to finish at shutdown.
; Increase this if you have very long running tasks.
stopwaitsecs = 6000

; When resorting to send SIGKILL to the program to terminate it
; send SIGKILL to its whole process group instead,
; taking care of its children as well.
killasgroup=true

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=997

;=============================================celery beat==================================================
[program:celerybeat]
; Set full path to celery program if using virtualenv
command=celery beat
; --schedule /var/lib/celery/beat.db --loglevel=INFO

; remove the -A myapp argument if you are not using an app instance

environment=PYTHONPATH=/home/rose/MStore

directory=/home/rose/MStore/scripts
user=rose
numprocs=1
stdout_logfile=/var/log/celery/beat.log
stderr_logfile=/var/log/celery/beat.log
autostart=false
autorestart=true
startsecs=10

; if rabbitmq is supervised, set its priority higher
; so it starts first
priority=996